# diane_web.service
[Unit]
Description=Diane Web Interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=diane
WorkingDirectory=/home/diane/diane
# Explicitly call the venv python on your web script:
ExecStart=/home/diane/diane/.venv/bin/python3 /home/diane/diane/diane_web.py
# Ensure PATH is correct so any imports work
Environment=PATH=/home/diane/diane/.venv/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
