# diane.service
[Unit]
Description=Diane Voice Assistant
After=network.target sound.target

[Service]
User=diane
WorkingDirectory=/opt/diane
ExecStart=/home/diane/diane_venv/bin/python3 /opt/diane/voice_llama_chat.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
