# voice_llama_chat.service
[Unit]
Description=Diane Voice Assistant (Precise wake-word)
After=network.target bluetooth.service

[Service]
Type=simple
User=diane
WorkingDirectory=/home/diane/diane
ExecStart=/home/diane/diane/voice_llama_chat.py \
  --llm-model "/mnt/ssd/models/L3.2-Rogue-Creative-Instruct-7B-D_AU-Q4_k_m.gguf" \
  --tts-model "/mnt/ssd/models/tts/en_US-amy-medium.onnx" \
  --precision-model "/mnt/ssd/models/hotword/precise_diane.pb" \
  --vad-mode "1" \
  --synonyms "/mnt/ssd/voice_config/synonyms.json" \
  --audio-input 1 \
  --audio-output 0 \
  --web-key "my_secret_key" \
  --log-level "INFO"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
